# deploy/ladder-dragon-log-export.service
# SPDX-License-Identifier: MIT
# Copyright (c) 2026 IURII Potekhin
# Purpose: define the systemd unit for ladder-dragon-log-export.
[Unit]
Description=Export sanitized Ladder Dragon logs for protected dashboard

[Service]
Type=oneshot
User=root
Group=www-data
Environment=BOT_LOG_EXPORT_DIR=/var/lib/ladder-dragon/logs
Environment=BOT_LOG_RETENTION_DAYS=7
Environment=BOT_LOG_MAX_BYTES=5242880
Environment=BOT_LOG_CURRENT_LINES=3000
ExecStart=/home/bot/apps/binance_bot/.venv/bin/python /home/bot/apps/binance_bot/deploy/export_sanitized_logs.py
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX
ReadWritePaths=/var/lib/ladder-dragon/logs
UMask=0027
