# deploy/ladder-dragon-backup.timer
# Copyright (c) 2026 IURII Potekhin / Ladder Dragon. All rights reserved.
# Назначение файла и опасные границы логики должны оставаться понятными при сопровождении.
[Unit]
Description=Daily Ladder Dragon private backup

[Timer]
OnCalendar=*-*-* 02:20:00
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
