# deploy/ladder-dragon-depth-archive.timer
# SPDX-License-Identifier: MIT
# Copyright (c) 2026 IURII Potekhin
# Purpose: sample public depth regularly across several volatility regimes.
[Unit]
Description=Schedule Ladder Dragon public depth archives

[Timer]
OnBootSec=20m
OnUnitActiveSec=1h
RandomizedDelaySec=5m
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
