# deploy/ladder-dragon-soak-audit.service
# SPDX-License-Identifier: MIT
# Copyright (c) 2026 IURII Potekhin
# Purpose: generate signed production soak evidence under root isolation.
[Unit]
Description=Ladder Dragon signed production soak audit
After=local-fs.target

[Service]
Type=oneshot
User=bot
Group=bot
Environment=PROJECT_DIR=/home/bot/apps/binance_bot
ExecStart=/usr/local/bin/ladder-dragon-soak-audit generate
ExecStartPost=+/usr/local/bin/ladder-dragon-soak-audit sign
Nice=15
IOSchedulingClass=idle
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
ReadOnlyPaths=-/run/mybot /home/bot/apps/binance_bot
ReadOnlyPaths=-/etc/ladder-dragon/telegram.env
ReadWritePaths=/var/lib/ladder-dragon/soak /etc/ladder-dragon
UMask=0007

[Install]
WantedBy=multi-user.target
