# deploy/ladder-dragon-backup.service
[Unit]
Description=Ladder Dragon private configuration and SQLite backup

[Service]
Type=oneshot
EnvironmentFile=/etc/ladder-dragon/backup.env
ExecStart=/home/bot/apps/binance_bot/deploy/backup_raspberry_pi.sh
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/home/bot/apps/binance_bot
# Динамические каталоги timestamp в backups и backups-public создаются во время
# запуска. Внешний диск подключается отдельным точечным bind в drop-in;
# общий /mnt здесь намеренно не пробрасывается, чтобы не скрывать вложенный mount.
ReadWritePaths=/var/lib/ladder-dragon
UMask=0077
