# deploy/pi-dashboard.service
# SPDX-License-Identifier: MIT
# Copyright (c) 2026 IURII Potekhin
# Purpose: define the systemd unit for pi-dashboard.
[Unit]
Description=Ladder Dragon private dashboard
After=network-online.target mybot.service

[Service]
Type=simple
User=bot
Group=bot
# Encrypted public backup metadata is root:www-data and remains read-only in
# this service's filesystem namespace.
SupplementaryGroups=www-data
WorkingDirectory=/home/bot/apps/binance_bot
EnvironmentFile=/home/bot/apps/binance_bot/.env.dashboard
ExecStart=/home/bot/apps/binance_bot/.venv/bin/python -m bin.run_dashboard
Restart=always
RestartSec=2
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectProc=invisible
ProtectSystem=strict
ProtectHome=read-only
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/home/bot/apps/binance_bot/FastAPI/pi-dashboard/data
ReadOnlyPaths=/home/bot/apps/binance_bot/db /home/bot/apps/binance_bot/logs -/run/mybot
ReadOnlyPaths=-/run/pi-watchdog/host-health.json
UMask=0077

[Install]
WantedBy=multi-user.target
