# deploy/pi-watchdog-v3.service
# SPDX-License-Identifier: MIT
# Copyright (c) 2026 IURII Potekhin
# Purpose: define the systemd unit for pi-watchdog-v3.
[Unit]
Description=Pi soft watchdog v3 (network + mybot)

[Service]
Type=oneshot
EnvironmentFile=-/etc/bot-alerts.env
Environment=WATCHDOG_LOG=/var/log/pi-watchdog/current.log
Environment=WATCHDOG_STATE=/run/pi-watchdog/state
ExecStart=/usr/local/bin/pi-watchdog_v3.sh
SuccessExitStatus=0
RuntimeDirectory=pi-watchdog
RuntimeDirectoryMode=0755
StateDirectory=pi-watchdog
StateDirectoryMode=0700
LogsDirectory=pi-watchdog
LogsDirectoryMode=0700
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectHostname=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
UMask=0077
