# deploy/pi-watchdog-v3.service
# SPDX-License-Identifier: MIT
# Copyright (c) 2026 IURII Potekhin
# Purpose: define the systemd unit for pi-watchdog-v3.
[Unit]
Description=Pi soft watchdog v3 (network + mybot)

[Service]
Type=oneshot
EnvironmentFile=-/etc/bot-alerts.env
ExecStart=/usr/local/bin/pi-watchdog_v3.sh
SuccessExitStatus=0
