# deploy/pi-watchdog-v3.service
# Copyright (c) 2026 IURII Potekhin / Ladder Dragon. All rights reserved.
# Назначение файла и опасные границы логики должны оставаться понятными при сопровождении.
[Unit]
Description=Pi soft watchdog v3 (network + mybot)

[Service]
Type=oneshot
EnvironmentFile=-/etc/bot-alerts.env
ExecStart=/usr/local/bin/pi-watchdog_v3.sh
SuccessExitStatus=0
